# var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (- installed on Thu Feb 15 03:46:07 2024)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
## Ispmanager acme.sh certs update task
37 3 * * *	/usr/local/mgr5/sbin/cron-core sbin/mgrctl -m core acmesh.certs.update >/dev/null 2>&1
## Ispmanager task which collect system statistics
*/5 * * * *	/usr/local/mgr5/sbin/cron-core sbin/mgrctl -m core sysinfostat >/dev/null 2>&1
## Ispmanager delete expire logs task
0 3 * * *	/usr/local/mgr5/sbin/cron-core sbin/mgrctl -m core deletearchivelogs >/dev/null 2>&1
## Ispmanager usagestat task
33 3 * * *	/usr/local/mgr5/sbin/cron-core sbin/usagestat >/dev/null 2>&1
## Ispmanager problems autosolve task (/usr/local/mgr5)
*/5 * * * *	/usr/local/mgr5/sbin/cron-core sbin/mgrctl -m core problems.autosolve >/dev/null 2>&1
## Ispmanager update package cache task
51 6 * * *	/usr/local/mgr5/sbin/cron-ispmgr sbin/pkgupdate.sh >/dev/null 2>&1
## Ispmanager stat handler (ispmgr)
36 11 * * *	/usr/local/mgr5/sbin/cron-ispmgr sbin/mgrctl -m ispmgr ispstat >/dev/null 2>&1
## Ispmanager periodic tasks handler (ispmgr)
* * * * *	/usr/local/mgr5/sbin/cron-ispmgr sbin/mgrctl -m ispmgr periodic >/dev/null 2>&1
## Ispmanager problems autosolve task (/usr/local/mgr5)
*/5 * * * *	/usr/local/mgr5/sbin/cron-ispmgr sbin/mgrctl -m ispmgr problems.autosolve >/dev/null 2>&1
## Ispmanager metric task
*/5 * * * *	/usr/local/mgr5/sbin/cron-ispmgr sbin/metric >/dev/null 2>&1
## Ispmanager upgrade task
6 5 * * *	/usr/local/mgr5/sbin/cron-ispmgr sbin/pkgupgrade.sh coremanager >/dev/null 2>&1
## Ispmanager letsencrypt cert update check
42 3 * * *	/usr/local/mgr5/sbin/cron-ispmgr sbin/mgrctl -m ispmgr letsencrypt.check.update >/dev/null 2>&1
## Check DS records from the domain registrar
4 4 * * 0	/usr/local/mgr5/sbin/cron-ispmgr sbin/mgrctl -m ispmgr checkdsrecords >/dev/null 2>&1
## Ispmanager services monitoring task (srvmon)
*/15 * * * *	/usr/local/mgr5/sbin/cron-ispmgr addon/srvmon >/dev/null 2>&1
