# var/lib/dpkg/info/isp-php74-fpm.list
/.
/etc
/etc/init.d
/etc/init.d/php-fpm74
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/php-fpm74.service
/opt
/opt/php74
/opt/php74/etc
/opt/php74/etc/php-fpm.conf
/opt/php74/etc/php-fpm.d
/opt/php74/etc/php-fpm.d/pool.d
/opt/php74/sbin
/opt/php74/sbin/php-fpm
/opt/php74/share
/opt/php74/share/man
/opt/php74/share/man/man8
/opt/php74/share/man/man8/php-fpm.8
/opt/php74/var
/opt/php74/var/run
/usr
/usr/share
/usr/share/doc
/usr/share/doc/isp-php74-fpm
/usr/share/doc/isp-php74-fpm/changelog.Debian.gz
