# var/lib/dpkg/info/isp-php56-fpm.conffiles
/opt/php56/etc/php-fpm.conf
/etc/init.d/php-fpm56
