# var/lib/dpkg/info/isp-php56-fpm.list
/.
/etc
/etc/init.d
/etc/init.d/php-fpm56
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/php-fpm56.service
/opt
/opt/php56
/opt/php56/etc
/opt/php56/etc/php-fpm.conf
/opt/php56/etc/php-fpm.d
/opt/php56/etc/php-fpm.d/pool.d
/opt/php56/sbin
/opt/php56/sbin/php-fpm
/opt/php56/share
/opt/php56/share/man
/opt/php56/share/man/man8
/opt/php56/share/man/man8/php-fpm.8
/opt/php56/var
/opt/php56/var/run
/usr
/usr/share
/usr/share/doc
/usr/share/doc/isp-php56-fpm
/usr/share/doc/isp-php56-fpm/changelog.Debian.gz
