# var/lib/dpkg/info/isp-php81-fpm.list
/.
/etc
/etc/init.d
/etc/init.d/php-fpm81
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/php-fpm81.service
/opt
/opt/php81
/opt/php81/etc
/opt/php81/etc/php-fpm.conf
/opt/php81/etc/php-fpm.d
/opt/php81/etc/php-fpm.d/pool.d
/opt/php81/sbin
/opt/php81/sbin/php-fpm
/opt/php81/share
/opt/php81/share/man
/opt/php81/share/man/man8
/opt/php81/share/man/man8/php-fpm.8
/opt/php81/var
/opt/php81/var/run
/usr
/usr/share
/usr/share/doc
/usr/share/doc/isp-php81-fpm
/usr/share/doc/isp-php81-fpm/changelog.Debian.gz
