# var/lib/dpkg/info/isp-php74-fpm.conffiles
/opt/php74/etc/php-fpm.conf
/etc/init.d/php-fpm74
