# var/lib/dpkg/info/ispmanager-pkg-httpd-itk.conffiles
/usr/local/mgr5/etc/conf.d/httpd-itk.conf
/etc/apache2/conf.d/letsencrypt.conf
/lib/systemd/system/apache2-isp@.service
/etc/apache2/conf.d/disabled.conf
/etc/apache2/conf.d/letsencrypt.conf
/etc/apache2/conf.d/z1_home.conf
