# var/lib/dpkg/info/ispmanager-pkg-httpd-itk.list
/.
/etc
/etc/apache2
/etc/apache2/conf.d
/etc/apache2/conf.d/disabled.conf
/etc/apache2/conf.d/letsencrypt.conf
/etc/apache2/conf.d/z1_home.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/apache2-isp@.service
/usr
/usr/local
/usr/local/mgr5
/usr/local/mgr5/etc
/usr/local/mgr5/etc/conf.d
/usr/local/mgr5/etc/conf.d/httpd-itk.conf
/usr/share
/usr/share/doc
/usr/share/doc/ispmanager-pkg-httpd-itk
/usr/share/doc/ispmanager-pkg-httpd-itk/changelog.Debian.gz
/usr/share/doc/ispmanager-pkg-httpd-itk/copyright
