# var/lib/dpkg/info/ispmanager-pkg-php.conffiles
/usr/local/mgr5/etc/conf.d/php.conf
/etc/cron.daily/phpsess_cron
/etc/cron.daily/phpsess_cron
