# var/lib/dpkg/info/postgresql-common.conffiles
/etc/apt/apt.conf.d/01autoremove-postgresql
/etc/init.d/postgresql
/etc/logrotate.d/postgresql-common
/etc/sysctl.d/30-postgresql-shm.conf
