# var/lib/dpkg/info/proftpd-basic.prerm
#!/bin/sh

set -e

invoke-rc.d proftpd stop || true


