# var/lib/dpkg/info/proftpd-basic.conffiles
/etc/default/proftpd
/etc/ftpusers
/etc/init.d/proftpd
/etc/logrotate.d/proftpd-basic
/etc/pam.d/proftpd
/etc/proftpd/blacklist.dat
/etc/proftpd/dhparams.pem
