# var/lib/dpkg/info/proftpd-basic.list
/.
/etc
/etc/cron.monthly
/etc/default
/etc/default/proftpd
/etc/ftpusers
/etc/init.d
/etc/init.d/proftpd
/etc/logrotate.d
/etc/logrotate.d/proftpd-basic
/etc/pam.d
/etc/pam.d/proftpd
/etc/proftpd
/etc/proftpd/blacklist.dat
/etc/proftpd/conf.d
/etc/proftpd/dhparams.pem
/usr
/usr/bin
/usr/bin/ftpcount
/usr/bin/ftpdctl
/usr/bin/ftptop
/usr/bin/ftpwho
/usr/lib
/usr/lib/proftpd
/usr/lib/proftpd/mod_auth_otp.so
/usr/lib/proftpd/mod_ban.so
/usr/lib/proftpd/mod_copy.so
/usr/lib/proftpd/mod_ctrls_admin.so
/usr/lib/proftpd/mod_deflate.so
/usr/lib/proftpd/mod_dynmasq.so
/usr/lib/proftpd/mod_exec.so
/usr/lib/proftpd/mod_facl.so
/usr/lib/proftpd/mod_ifsession.so
/usr/lib/proftpd/mod_ifversion.so
/usr/lib/proftpd/mod_load.so
/usr/lib/proftpd/mod_memcache.so
/usr/lib/proftpd/mod_quotatab.so
/usr/lib/proftpd/mod_quotatab_file.so
/usr/lib/proftpd/mod_quotatab_radius.so
/usr/lib/proftpd/mod_quotatab_sql.so
/usr/lib/proftpd/mod_radius.so
/usr/lib/proftpd/mod_ratio.so
/usr/lib/proftpd/mod_readme.so
/usr/lib/proftpd/mod_redis.so
/usr/lib/proftpd/mod_rewrite.so
/usr/lib/proftpd/mod_sftp.so
/usr/lib/proftpd/mod_sftp_pam.so
/usr/lib/proftpd/mod_sftp_sql.so
/usr/lib/proftpd/mod_shaper.so
/usr/lib/proftpd/mod_site_misc.so
/usr/lib/proftpd/mod_sql.so
/usr/lib/proftpd/mod_sql_passwd.so
/usr/lib/proftpd/mod_tls.so
/usr/lib/proftpd/mod_tls_memcache.so
/usr/lib/proftpd/mod_tls_redis.so
/usr/lib/proftpd/mod_unique_id.so
/usr/lib/proftpd/mod_wrap.so
/usr/lib/proftpd/mod_wrap2.so
/usr/lib/proftpd/mod_wrap2_file.so
/usr/lib/proftpd/mod_wrap2_redis.so
/usr/lib/proftpd/mod_wrap2_sql.so
/usr/sbin
/usr/sbin/auth-otp
/usr/sbin/ftpasswd
/usr/sbin/ftpquota
/usr/sbin/ftpscrub
/usr/sbin/ftpshut
/usr/sbin/ftpstats
/usr/sbin/proftpd
/usr/sbin/proftpd-gencert
/usr/share
/usr/share/doc
/usr/share/doc/proftpd-basic
/usr/share/doc/proftpd-basic/NEWS.Debian.gz
/usr/share/doc/proftpd-basic/README.Debian.gz
/usr/share/doc/proftpd-basic/TODO.Debian
/usr/share/doc/proftpd-basic/changelog.Debian.gz
/usr/share/doc/proftpd-basic/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/proftpd-basic
/usr/share/locale
/usr/share/locale/bg_BG
/usr/share/locale/bg_BG/LC_MESSAGES
/usr/share/locale/bg_BG/LC_MESSAGES/proftpd.mo
/usr/share/locale/en_US
/usr/share/locale/en_US/LC_MESSAGES
/usr/share/locale/en_US/LC_MESSAGES/proftpd.mo
/usr/share/locale/es_ES
/usr/share/locale/es_ES/LC_MESSAGES
/usr/share/locale/es_ES/LC_MESSAGES/proftpd.mo
/usr/share/locale/fr_FR
/usr/share/locale/fr_FR/LC_MESSAGES
/usr/share/locale/fr_FR/LC_MESSAGES/proftpd.mo
/usr/share/locale/it_IT
/usr/share/locale/it_IT/LC_MESSAGES
/usr/share/locale/it_IT/LC_MESSAGES/proftpd.mo
/usr/share/locale/ja_JP
/usr/share/locale/ja_JP/LC_MESSAGES
/usr/share/locale/ja_JP/LC_MESSAGES/proftpd.mo
/usr/share/locale/ko_KR
/usr/share/locale/ko_KR/LC_MESSAGES
/usr/share/locale/ko_KR/LC_MESSAGES/proftpd.mo
/usr/share/locale/ru_RU
/usr/share/locale/ru_RU/LC_MESSAGES
/usr/share/locale/ru_RU/LC_MESSAGES/proftpd.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/proftpd.mo
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/proftpd.mo
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ftpcount.1.gz
/usr/share/man/man1/ftpmail.1.gz
/usr/share/man/man1/ftpquota.1.gz
/usr/share/man/man1/ftptop.1.gz
/usr/share/man/man1/ftpwho.1.gz
/usr/share/man/man5
/usr/share/man/man5/ftpusers.5.gz
/usr/share/man/man5/proftpd.conf.5.gz
/usr/share/man/man5/xferlog.5.gz
/usr/share/man/man8
/usr/share/man/man8/auth-otp.8.gz
/usr/share/man/man8/ftpasswd.8.gz
/usr/share/man/man8/ftpdctl.8.gz
/usr/share/man/man8/ftpscrub.8.gz
/usr/share/man/man8/ftpshut.8.gz
/usr/share/man/man8/ftpstats.8.gz
/usr/share/man/man8/proftpd.8.gz
/usr/share/proftpd
/usr/share/proftpd/templates
/usr/share/proftpd/templates/geoip.conf
/usr/share/proftpd/templates/ldap.conf
/usr/share/proftpd/templates/mod_snmp.conf
/usr/share/proftpd/templates/modules.conf
/usr/share/proftpd/templates/proftpd.conf
/usr/share/proftpd/templates/sql.conf
/usr/share/proftpd/templates/tls.conf
/usr/share/proftpd/templates/virtuals.conf
/usr/share/proftpd/templates/welcome.msg
/var
/var/log
/var/log/proftpd
/usr/sbin/in.proftpd
