# var/lib/dpkg/info/qemu-guest-agent.list
/.
/etc
/etc/init.d
/etc/init.d/qemu-guest-agent
/etc/qemu
/etc/qemu/fsfreeze-hook
/etc/qemu/fsfreeze-hook.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/qemu-guest-agent.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-qemu-guest-agent.rules
/usr
/usr/sbin
/usr/sbin/qemu-ga
/usr/share
/usr/share/doc
/usr/share/doc/qemu-guest-agent
/usr/share/doc/qemu-guest-agent/changelog.Debian.gz
/usr/share/doc/qemu-guest-agent/copyright
/usr/share/doc/qemu-guest-agent/qapi-schema.json.gz
/usr/share/doc/qemu-guest-agent/qemu-ga-ref.html
/usr/share/doc/qemu-guest-agent/qemu-ga-ref.txt.gz
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/qemu-ga-ref.7.gz
/usr/share/man/man8
/usr/share/man/man8/qemu-ga.8.gz
