# var/lib/dpkg/info/tuned.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/com.redhat.tuned.conf
/etc/grub.d
/etc/grub.d/00_tuned
/etc/modprobe.d
/etc/modprobe.d/tuned.conf
/etc/tuned
/etc/tuned/active_profile
/etc/tuned/bootcmdline
/etc/tuned/cpu-partitioning-variables.conf
/etc/tuned/profile_mode
/etc/tuned/realtime-variables.conf
/etc/tuned/realtime-virtual-guest-variables.conf
/etc/tuned/realtime-virtual-host-variables.conf
/etc/tuned/sap-hana-vmware-variables.conf
/etc/tuned/tuned-main.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/tuned.service
/usr
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/tuned
/usr/lib/python3/dist-packages/tuned/__init__.py
/usr/lib/python3/dist-packages/tuned/admin
/usr/lib/python3/dist-packages/tuned/admin/__init__.py
/usr/lib/python3/dist-packages/tuned/admin/admin.py
/usr/lib/python3/dist-packages/tuned/admin/dbus_controller.py
/usr/lib/python3/dist-packages/tuned/admin/exceptions.py
/usr/lib/python3/dist-packages/tuned/consts.py
/usr/lib/python3/dist-packages/tuned/daemon
/usr/lib/python3/dist-packages/tuned/daemon/__init__.py
/usr/lib/python3/dist-packages/tuned/daemon/application.py
/usr/lib/python3/dist-packages/tuned/daemon/controller.py
/usr/lib/python3/dist-packages/tuned/daemon/daemon.py
/usr/lib/python3/dist-packages/tuned/exceptions.py
/usr/lib/python3/dist-packages/tuned/exports
/usr/lib/python3/dist-packages/tuned/exports/__init__.py
/usr/lib/python3/dist-packages/tuned/exports/controller.py
/usr/lib/python3/dist-packages/tuned/exports/dbus_exporter.py
/usr/lib/python3/dist-packages/tuned/exports/interfaces.py
/usr/lib/python3/dist-packages/tuned/hardware
/usr/lib/python3/dist-packages/tuned/hardware/__init__.py
/usr/lib/python3/dist-packages/tuned/hardware/device_matcher.py
/usr/lib/python3/dist-packages/tuned/hardware/device_matcher_udev.py
/usr/lib/python3/dist-packages/tuned/hardware/inventory.py
/usr/lib/python3/dist-packages/tuned/logs.py
/usr/lib/python3/dist-packages/tuned/monitors
/usr/lib/python3/dist-packages/tuned/monitors/__init__.py
/usr/lib/python3/dist-packages/tuned/monitors/base.py
/usr/lib/python3/dist-packages/tuned/monitors/monitor_disk.py
/usr/lib/python3/dist-packages/tuned/monitors/monitor_load.py
/usr/lib/python3/dist-packages/tuned/monitors/monitor_net.py
/usr/lib/python3/dist-packages/tuned/monitors/repository.py
/usr/lib/python3/dist-packages/tuned/patterns.py
/usr/lib/python3/dist-packages/tuned/plugins
/usr/lib/python3/dist-packages/tuned/plugins/__init__.py
/usr/lib/python3/dist-packages/tuned/plugins/base.py
/usr/lib/python3/dist-packages/tuned/plugins/decorators.py
/usr/lib/python3/dist-packages/tuned/plugins/exceptions.py
/usr/lib/python3/dist-packages/tuned/plugins/hotplug.py
/usr/lib/python3/dist-packages/tuned/plugins/instance
/usr/lib/python3/dist-packages/tuned/plugins/instance/__init__.py
/usr/lib/python3/dist-packages/tuned/plugins/instance/factory.py
/usr/lib/python3/dist-packages/tuned/plugins/instance/instance.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_audio.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_bootloader.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_cpu.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_disk.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_eeepc_she.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_modules.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_mounts.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_net.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_script.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_scsi_host.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_selinux.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_sysctl.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_sysfs.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_systemd.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_usb.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_video.py
/usr/lib/python3/dist-packages/tuned/plugins/plugin_vm.py
/usr/lib/python3/dist-packages/tuned/plugins/repository.py
/usr/lib/python3/dist-packages/tuned/profiles
/usr/lib/python3/dist-packages/tuned/profiles/__init__.py
/usr/lib/python3/dist-packages/tuned/profiles/exceptions.py
/usr/lib/python3/dist-packages/tuned/profiles/factory.py
/usr/lib/python3/dist-packages/tuned/profiles/functions
/usr/lib/python3/dist-packages/tuned/profiles/functions/__init__.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/base.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_assertion.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_assertion_non_equal.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_cpulist2hex.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_cpulist2hex_invert.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_cpulist_invert.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_cpulist_online.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_cpulist_pack.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_cpulist_present.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_cpulist_unpack.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_exec.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_hex2cpulist.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_kb2s.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_s2kb.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_strip.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/function_virt_check.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/functions.py
/usr/lib/python3/dist-packages/tuned/profiles/functions/repository.py
/usr/lib/python3/dist-packages/tuned/profiles/loader.py
/usr/lib/python3/dist-packages/tuned/profiles/locator.py
/usr/lib/python3/dist-packages/tuned/profiles/merger.py
/usr/lib/python3/dist-packages/tuned/profiles/profile.py
/usr/lib/python3/dist-packages/tuned/profiles/unit.py
/usr/lib/python3/dist-packages/tuned/profiles/variables.py
/usr/lib/python3/dist-packages/tuned/storage
/usr/lib/python3/dist-packages/tuned/storage/__init__.py
/usr/lib/python3/dist-packages/tuned/storage/factory.py
/usr/lib/python3/dist-packages/tuned/storage/interfaces.py
/usr/lib/python3/dist-packages/tuned/storage/pickle_provider.py
/usr/lib/python3/dist-packages/tuned/storage/storage.py
/usr/lib/python3/dist-packages/tuned/units
/usr/lib/python3/dist-packages/tuned/units/__init__.py
/usr/lib/python3/dist-packages/tuned/units/manager.py
/usr/lib/python3/dist-packages/tuned/utils
/usr/lib/python3/dist-packages/tuned/utils/__init__.py
/usr/lib/python3/dist-packages/tuned/utils/commands.py
/usr/lib/python3/dist-packages/tuned/utils/global_config.py
/usr/lib/python3/dist-packages/tuned/utils/nettool.py
/usr/lib/python3/dist-packages/tuned/utils/plugin_loader.py
/usr/lib/python3/dist-packages/tuned/utils/polkit.py
/usr/lib/python3/dist-packages/tuned/version.py
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/tuned.conf
/usr/lib/tuned
/usr/lib/tuned/atomic-guest
/usr/lib/tuned/atomic-guest/tuned.conf
/usr/lib/tuned/atomic-host
/usr/lib/tuned/atomic-host/tuned.conf
/usr/lib/tuned/balanced
/usr/lib/tuned/balanced/tuned.conf
/usr/lib/tuned/cpu-partitioning
/usr/lib/tuned/cpu-partitioning/00-tuned-pre-udev.sh
/usr/lib/tuned/cpu-partitioning/script.sh
/usr/lib/tuned/cpu-partitioning/tuned.conf
/usr/lib/tuned/default
/usr/lib/tuned/default/tuned.conf
/usr/lib/tuned/defirqaffinity.py
/usr/lib/tuned/desktop
/usr/lib/tuned/desktop/tuned.conf
/usr/lib/tuned/desktop-powersave
/usr/lib/tuned/desktop-powersave/tuned.conf
/usr/lib/tuned/enterprise-storage
/usr/lib/tuned/enterprise-storage/tuned.conf
/usr/lib/tuned/functions
/usr/lib/tuned/laptop-ac-powersave
/usr/lib/tuned/laptop-ac-powersave/script.sh
/usr/lib/tuned/laptop-ac-powersave/tuned.conf
/usr/lib/tuned/laptop-battery-powersave
/usr/lib/tuned/laptop-battery-powersave/tuned.conf
/usr/lib/tuned/latency-performance
/usr/lib/tuned/latency-performance/tuned.conf
/usr/lib/tuned/mssql
/usr/lib/tuned/mssql/tuned.conf
/usr/lib/tuned/network-latency
/usr/lib/tuned/network-latency/tuned.conf
/usr/lib/tuned/network-throughput
/usr/lib/tuned/network-throughput/tuned.conf
/usr/lib/tuned/oracle
/usr/lib/tuned/oracle/tuned.conf
/usr/lib/tuned/powersave
/usr/lib/tuned/powersave/script.sh
/usr/lib/tuned/powersave/tuned.conf
/usr/lib/tuned/realtime
/usr/lib/tuned/realtime/script.sh
/usr/lib/tuned/realtime/tuned.conf
/usr/lib/tuned/realtime-virtual-guest
/usr/lib/tuned/realtime-virtual-guest/tuned.conf
/usr/lib/tuned/realtime-virtual-host
/usr/lib/tuned/realtime-virtual-host/find-lapictscdeadline-optimal.sh
/usr/lib/tuned/realtime-virtual-host/script.sh
/usr/lib/tuned/realtime-virtual-host/tuned.conf
/usr/lib/tuned/recommend.d
/usr/lib/tuned/recommend.d/50-tuned.conf
/usr/lib/tuned/sap-hana
/usr/lib/tuned/sap-hana/tuned.conf
/usr/lib/tuned/sap-hana-vmware
/usr/lib/tuned/sap-hana-vmware/tuned.conf
/usr/lib/tuned/sap-netweaver
/usr/lib/tuned/sap-netweaver/script.sh
/usr/lib/tuned/sap-netweaver/tuned.conf
/usr/lib/tuned/server-powersave
/usr/lib/tuned/server-powersave/tuned.conf
/usr/lib/tuned/spindown-disk
/usr/lib/tuned/spindown-disk/script.sh
/usr/lib/tuned/spindown-disk/tuned.conf
/usr/lib/tuned/throughput-performance
/usr/lib/tuned/throughput-performance/tuned.conf
/usr/lib/tuned/virtual-guest
/usr/lib/tuned/virtual-guest/tuned.conf
/usr/lib/tuned/virtual-host
/usr/lib/tuned/virtual-host/tuned.conf
/usr/sbin
/usr/sbin/tuned
/usr/sbin/tuned-adm
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/tuned-adm
/usr/share/doc
/usr/share/doc/tuned
/usr/share/doc/tuned/README
/usr/share/doc/tuned/README.NFV
/usr/share/doc/tuned/TIPS.txt
/usr/share/doc/tuned/changelog.Debian.gz
/usr/share/doc/tuned/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/tuned-main.conf.5.gz
/usr/share/man/man5/tuned.conf.5.gz
/usr/share/man/man7
/usr/share/man/man7/tuned-profiles-atomic.7.gz
/usr/share/man/man7/tuned-profiles-compat.7.gz
/usr/share/man/man7/tuned-profiles-cpu-partitioning.7.gz
/usr/share/man/man7/tuned-profiles-mssql.7.gz
/usr/share/man/man7/tuned-profiles-nfv-guest.7.gz
/usr/share/man/man7/tuned-profiles-nfv-host.7.gz
/usr/share/man/man7/tuned-profiles-oracle.7.gz
/usr/share/man/man7/tuned-profiles-realtime.7.gz
/usr/share/man/man7/tuned-profiles-sap-hana.7.gz
/usr/share/man/man7/tuned-profiles-sap.7.gz
/usr/share/man/man7/tuned-profiles.7.gz
/usr/share/man/man8
/usr/share/man/man8/tuned-adm.8.gz
/usr/share/man/man8/tuned.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/com.redhat.tuned.policy
/usr/share/python3
/usr/share/python3/runtime.d
/usr/share/python3/runtime.d/tuned.rtupdate
